# roles/consul/maintenance-snapshot/templates/consul-snapshot.timer.j2
[Unit]
Description=Run scheduled Consul Raft snapshot

[Timer]
# Ref: https://www.freedesktop.org/software/systemd/man/systemd.time.html#Calendar%20Events
OnCalendar={{ consul_snapshot_timing }}
Persistent=true

[Install]
WantedBy=timers.target
